# root/etc/services.d/plex-setup/run
#!/usr/bin/with-contenv bash
python run.py
s6-pause th
